# .devcontainer/build/shell/profile.d/alias.sh
if type code-insiders > /dev/null 2>&1 && ! type code > /dev/null 2>&1; then
    alias code=code-insiders
fi

alias ll='ls -l -a'
